# config/reddit_bot_supervisor.conf
[program:uwsgi]
command=make run-bot
directory=/home/ubuntu/EmojiFi/
autostart=true
autorestart=true
user=ubuntu
