# config/uwsgi_supervisord.conf
[program:uwsgi]
command=sudo uwsgi --emperor /etc/uwsgi/vassals --uid www-data --gid www-data
directory=/
autostart=true
autorestart=true
priority=50
user=ubuntu
